# start_scripts.service
[Unit]
Description=Start Python applications via screen sessions
After=network.target

[Service]
Type=simple
User=scad-pi
WorkingDirectory=/home/scad-pi/raspberry-pi-server
ExecStart=/bin/bash /home/scad-pi/raspberry-pi-server/start_scripts.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target